# os/stratusd/stratusd.service
[Unit]
Description=The stratus game streaming daemon
After=network.target
Requires=stratusd-sway.service
After=stratusd-sway.service

[Service]
Environment="WAYLAND_DISPLAY=wayland-1"
Environment="DISPLAY=:0"
EnvironmentFile=/etc/stratusd.conf
ExecStartPre=/bin/sh -c "until ping -c 1 $(echo $STRATUSD_BACKEND_URL | awk -F[/:] '{print $4}'); do sleep 1; done"
ExecStart=/usr/bin/stratusd --daemon
Restart=always

[Install]
WantedBy=default.target
